# device/setup/kakabox-reset.service
# Hintergrund-Watcher: Roter Knopf 10s halten → WLAN vergessen + Reboot.
#
# Läuft IMMER, unabhängig vom WLAN-Status oder der Hauptanwendung.
[Unit]
Description=Kakabox WiFi reset watcher (red button hold)
After=multi-user.target

[Service]
# Root nötig: nmcli connection delete (System-Connections) + systemctl reboot
Type=simple
User=root
ExecStart=/home/riffi/Dokumente/kakabox/device/.venv/bin/python /home/riffi/Dokumente/kakabox/device/setup/kakabox-reset-watcher.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
