# device/setup/kakabox.service
# systemd-Unit für die Kakabox-Hauptanwendung.
#
# Startet beim Boot (enabled), unabhängig vom WLAN — die Box soll auch ohne
# Internet mit den lokal gecachten Kakas funktionieren. Backend-Funktionen
# (Sync, Heartbeat) sind innerhalb von main.py best-effort und legen lahm,
# wenn kein WLAN da ist.
#
# Manuell:
#   sudo systemctl start kakabox.service
#   sudo journalctl -u kakabox -f
[Unit]
Description=Kakabox device main loop
After=sound.target
# StartLimit-Fenster auf Unit-Ebene: höchstens 10 Restart-Versuche in 2 Min,
# danach lässt systemd den Service "failed". Verhindert Endlos-Restart-Loops,
# falls main.py beim Start sauber crasht (z.B. kaputte config.json).
StartLimitBurst=10
StartLimitIntervalSec=120

[Service]
Type=simple
User=riffi
Group=riffi
WorkingDirectory=/home/riffi/Dokumente/kakabox/device
ExecStart=/home/riffi/Dokumente/kakabox/device/.venv/bin/python /home/riffi/Dokumente/kakabox/device/main.py
# Always: auch bei sauberem Exit neu starten (z.B. unerwartetes return aus
# main()). systemd supervised PID-event-driven — kein Polling, keine
# Performance-Kosten. Manuelles 'systemctl stop' wird respektiert.
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# ──────────────────────────────────────────────────────────────────────────
# Hardening — Kindergerät, deshalb so eng wie möglich.
# Sollte irgendwann jemand RCE in main.py finden, soll der Schaden begrenzt
# bleiben. Jede Option einzeln getestet — die unten erwähnten brechen
# Hardware/Audio und sind absichtlich nicht gesetzt.
# ──────────────────────────────────────────────────────────────────────────
# /usr, /boot, /etc read-only
ProtectSystem=full
# andere Home-Dirs nicht zugreifbar
ProtectHome=read-only
# einziger r/w-Pfad: das eigene Device-Verzeichnis (config.json, audio_cache, …)
ReadWritePaths=/home/riffi/Dokumente/kakabox/device
# eigener /tmp-Namespace, nichts persistiert über Restart
PrivateTmp=true
# cgroups read-only — kein Tampering mit Resource-Limits
ProtectControlGroups=true
# /proc/sys/kernel/hostname read-only
ProtectHostname=true
# Systemzeit unveränderbar
ProtectClock=true
# /dev/kmsg gesperrt — keine Spurenverwischung
ProtectKernelLogs=true
# personality(2) eingefroren
LockPersonality=true
# keine neuen Namespaces — kein Container-Escape-Trick
RestrictNamespaces=true
# kein Realtime-Scheduling, sonst lässt sich der Scheduler quälen
RestrictRealtime=true
# neu erstellte Dateien nur für riffi lesbar (config.json enthält API-Token)
UMask=0077

# Bewusst NICHT gesetzt: NoNewPrivileges=true
#   → würde den sudo-Aufruf für /usr/local/bin/kakabox-wifi-nuke
#     (Roter Knopf 10s halten = WLAN-Reset) blockieren. Ersatz wäre ein
#     root-Daemon hinter einer Pipe — aktuell nicht gerechtfertigter
#     Mehrwert, weil der sudoers-Drop-in eng auf genau diesen einen Pfad
#     ohne Argumente eingegrenzt ist (siehe sudoers-kakabox).

[Install]
WantedBy=multi-user.target
